# docker/bisq/config/s6-overlay/s6-rc.d/vnc/run
#!/command/with-contenv bash
exec s6-setuidgid bisq /home/bisq/start-vnc.sh