# docker/bisq/config/s6-overlay/s6-rc.d/novnc/run
#!/command/with-contenv bash
cd /usr/share/novnc
exec s6-setuidgid bisq websockify --web /usr/share/novnc 6080 localhost:5901