# docker/bisq/config/s6-overlay/s6-rc.d/resize-monitor/run
#!/command/with-contenv bash
exec s6-setuidgid bisq /home/bisq/monitor-resize.sh